# rootfs/etc/supervisor.d/mariadb.ini
[program:mariadb]
user=mysql
command=/usr/bin/mysqld 
numprocs=1
autostart=true
autorestart=true
