# rootfs/etc/supervisor.d/mariadb.ini
[program:mariadb]
user=%(ENV_APP_USER)s
command=%(ENV_MYSQLD_CMD)s
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
