# rootfs/etc/supervisor.d/mariadb.ini
[program:mariadb]
user=%(ENV_APP_USER)s
command=/usr/bin/mysqld
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisord/mariadb.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
redirect_stderr=true
